# extra-files.template/etc/uci-defaults/99_zzzz-post_boot.sh
#!/bin/sh
PATH=/usr/bin:/usr/sbin:/bin:/sbin
export PATH
cd /custom
lua post_boot.lua
rc=$?
exit $rc
