# extra-files.template/etc/hotplug.d/iface/99-ifup-vpn
#!/bin/sh
[ "$ACTION" = "ifup" -a "$(echo $INTERFACE | grep -i vpn)" != "" ] && {
    sysctl -ew net.ipv4.conf.${INTERFACE}.rp_filter=2
}
exit 0